# dolphin/launchd/com.dolphin.weekly-report.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.dolphin.weekly-report</string>

    <key>ProgramArguments</key>
    <array>
        <string>/Library/Frameworks/Python.framework/Versions/3.13/bin/python3</string>
        <string>/Users/johnasbury/Reachh/dolphin/reporting.py</string>
    </array>

    <key>WorkingDirectory</key>
    <string>/Users/johnasbury/Reachh/dolphin</string>

    <!-- Run every Sunday at 10 AM -->
    <key>StartCalendarInterval</key>
    <dict>
        <key>Weekday</key>
        <integer>0</integer>
        <key>Hour</key>
        <integer>10</integer>
        <key>Minute</key>
        <integer>0</integer>
    </dict>

    <key>StandardOutPath</key>
    <string>/Users/johnasbury/Reachh/dolphin/logs/weekly-report.stdout</string>
    <key>StandardErrorPath</key>
    <string>/Users/johnasbury/Reachh/dolphin/logs/weekly-report.stderr</string>
</dict>
</plist>
